# PrivilegedHelpers/cmux-sudo-helper/com.cmuxterm.sudo-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.cmuxterm.sudo-helper</string>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>@CMUX_BUNDLE_IDENTIFIER@</string>
	</array>
	<key>BundleProgram</key>
	<string>Contents/Library/LaunchServices/com.cmuxterm.sudo-helper</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/Logs/cmux-sudo-helper.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/Logs/cmux-sudo-helper.log</string>
</dict>
</plist>
