# configs/systemd/cmux-websockify.service
[Unit]
Description=Cmux websockify proxy
After=cmux-tigervnc.service
Requires=cmux-tigervnc.service

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/usr/bin/websockify --web=/usr/share/novnc 0.0.0.0:39380 127.0.0.1:5901
Restart=always
RestartSec=3
StandardOutput=append:/var/log/cmux/websockify.log
StandardError=append:/var/log/cmux/websockify.log

[Install]
WantedBy=cmux.target
